# .config/autostart/[code] code vs code.desktop
[Desktop Entry]
Exec=/home/dex/scripts/autostart/code_wsR.sh
Name=[code] code vscode
Type=Application
Version=1.0
