# .config/autostart/[Launch] Artha - dictionary.desktop
[Desktop Entry]
Exec=artha # just launch you got rulles
Name=[Launch] Artha - dictionary
Type=Application
Version=1.0
