# .config/autostart/[Panel] xfce.desktop
[Desktop Entry]
Exec=xfce4-panel
Name=[Panel] xfce
Type=Application
Version=1.0
