# .config/autostart/[Password Manager] Keepassxc.desktop
[Desktop Entry]
Exec=keepassxc
Name=[Password Manager] Keepassxc
Type=Application
Version=1.0
