# .config/autostart/[remote control] Deskflow.desktop
[Desktop Entry]
Exec=deskflow
Name=[remote control] Deskflow
Type=Application
Version=1.0
