# .config/autostart/[ssh github] keychain for github.desktop
[Desktop Entry]
Exec=keychain --eval ~/.ssh/githubarch
Name=[ssh github] keychain for github
Type=Application
Version=1.0
