# .config/autostart/[Virt-manager] got ram.desktop
[Desktop Entry]
Exec=sh -c \\\\"env SPICE_NOGRAB=1 virt-manager\\\\"
Name=[Virt-manager (no grab)] got ram
Type=Application
Version=1.0
