# .config/autostart/Devilspie2.desktop
[Desktop Entry]
Exec=devilspie2
Name=Devilspie2
Type=Application
Version=1.0
