# .config/autostart/run flameshot at startup.desktop
[Desktop Entry]
Exec=flameshot &
Name=[Program] flameshot [ # for screenshtos ]
Type=Application
Version=1.0
